# autoupdate.service
[Unit]
Description=Script de Atualização Automática

[Service]
Type=oneshot
ExcStart=/opt/htepil/autoupdate.sh

[Install]
WantedBy=multi-user.target
